# src/packaging/Debian/openafs-client.conffiles
/etc/openafs/afs.conf
/etc/init.d/openafs-client
